# keyleds.service
[Unit]
Description=Keyleadss deamon

[Service]
ExecStart=/usr/bin/keyledsd

[Install]
WantedBy=default.target